# packaging/systemd/hearth.service
[Unit]
Description=Hearth Personal Reticulum Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hearth
Group=hearth
WorkingDirectory=/opt/hearth
Environment=PYTHONUNBUFFERED=1
Environment=HEARTH_CONFIG=/etc/hearth/hearth.toml
ExecStart=/opt/hearth/.venv/bin/hearth-api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
